# features/isolation/data/etc/udev/rules.d/10-firmware-isolate.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/sbin/firmware-isolate"
